# target/linux/ar71xx/base-files/etc/uci-defaults/99_sun-enable-wifi
#!/bin/sh

. /lib/ar71xx.sh

board=$(ar71xx_board_name)

if [ "${board}" == "zsun-sdreader" ]; then
  uci set wireless.@wifi-device[0].disabled=0
  uci commit wireless
  wifi
fi

exit 0
